# templates/openclaw-gateway.service.tmpl
[Unit]
Description=OpenClaw Gateway
After=network-online.target
Wants=network-online.target

[Service]
ExecStart="/usr/bin/node" "/usr/lib/node_modules/openclaw/dist/index.js" gateway --port ${GATEWAY_PORT}
Restart=always
RestartSec=5
KillMode=process
Environment=NVIDIA_API_KEY=${NVIDIA_API_KEY}
Environment="MEM0_API_KEY=${MEM0_API_KEY}"
Environment=HOME=/home/${BOT_USER}
Environment="PATH=/home/${BOT_USER}/.local/bin:/home/${BOT_USER}/.npm-global/bin:/home/${BOT_USER}/bin:/usr/local/bin:/usr/bin:/bin"
Environment=OPENCLAW_GATEWAY_PORT=${GATEWAY_PORT}
Environment=OPENCLAW_GATEWAY_TOKEN=${GATEWAY_TOKEN}
Environment="OPENCLAW_SYSTEMD_UNIT=openclaw-gateway.service"
Environment=OPENCLAW_SERVICE_MARKER=openclaw
Environment=OPENCLAW_SERVICE_KIND=gateway

[Install]
WantedBy=default.target
